# launchd/jarvis.scheduler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>jarvis.scheduler</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/hanga/Desktop/claude/Jarvis/scheduler.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/hanga/Desktop/claude/Jarvis</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <!-- Run every 5 minutes to check backlog -->
    <key>StartInterval</key>
    <integer>300</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/hanga/Desktop/claude/Jarvis/logs/scheduler.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/hanga/Desktop/claude/Jarvis/logs/scheduler.error.log</string>
</dict>
</plist>
