# launchd/jarvis.watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>jarvis.watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/hanga/venv/bin/python3.14</string>
        <string>/Users/hanga/Desktop/claude/Jarvis/systems/file_watcher.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/hanga/Desktop/claude/Jarvis</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/hanga/venv/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/hanga/Desktop/claude/Jarvis</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/hanga/Desktop/claude/Jarvis/logs/watcher_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/hanga/Desktop/claude/Jarvis/logs/watcher_stderr.log</string>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
